# etc/supervisor/maglica_copy_worker.conf
[program:maglica_copy_worker]
command=/usr/bin/maglica_copy_worker
stdout_logfile=/var/log/maglica/copy_worker.log
stderr_logfile=/var/log/maglica/copy_worker.log