# etc/supervisor/maglica_host_worker.conf
[program:maglica_host_worker]
command=/usr/bin/maglica_host_worker
stdout_logfile=/var/log/maglica/host_worker.log
stderr_logfile=/var/log/maglica/host_worker.log