# etc/supervisor/maglica_client_worker.conf
[program:maglica_client_worker]
command=/usr/bin/maglica_client_worker
stdout_logfile=/var/log/maglica/client_worker.log
stderr_logfile=/var/log/maglica/client_worker.log